# services/read-sensors.socket
[Unit]
Description=FlatSat2 Read Sensors socket

[Socket]
ListenDatagram=/run/read-sensors/read-sensors.sock
SocketMode=0666

[Install]
WantedBy=sockets.target

